# submit_report.timer
[Unit]
Description = Submit timer

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
